# container/etc/cron.d/default
* * * * * /usr/local/bin/php /usr/local/empty.php
