# container/etc/cron.d/magento-crons
PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
* * * * * /bin/bash /var/www/html/scheduler_cron.sh --mode always
* * * * * /bin/bash /var/www/html/scheduler_cron.sh --mode default
