# container/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:crond]
command=cron